# case-studies/remcos/rules/sigma/remcos-registry-run-keys.yml
title: Remcos RAT Persistence Through Windows Run Registry Keys
id: bd547d52-df32-43f5-b3a8-b4eb8c38683d
status: experimental
description: Detects known Remcos RAT file and value names being added to Windows Registry Run keys for automatic startup.
author: James Banday
date: 2026-06-19
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: registry_set
    product: windows
detection:
    # Run and RunOnce keys can start a program automatically when a user signs in.
    selection_run_key:
        TargetObject|contains:
            - '\Software\Microsoft\Windows\CurrentVersion\Run\'
            - '\Software\Microsoft\Windows\CurrentVersion\RunOnce\'

    # This value name was identified during the completed Remcos investigation.
    selection_remcos_value:
        TargetObject|endswith: '\eastvillageeatery'

    # These filenames are known indicators from the same investigation.
    selection_remcos_file:
        Details|contains:
            - 'eastvillageeatery.exe'
            - '7hr8ftfxm.exe'

    # Alert only when a Run key change also contains a known Remcos indicator.
    condition: selection_run_key and 1 of selection_remcos_*
falsepositives:
    - A legitimate file could reuse one of these names; confirm the file hash and surrounding endpoint activity.
level: high
